# .config/supervisor/uwsgi.conf
# docker run 실행시 곧바로 uWSGI에 의해서 서버가 작동
[program:uwsgi]
command = /root/.pyenv/versions/used-book-store/bin/uwsgi --ini /srv/used-book-store/.config/uwsgi/uwsgi-app.ini